# vm-agent/vm-browser-agent.service
[Unit]
Description=VM Browser Agent - OAuth Proxy for CLI Tools
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/vm-agent/vm-browser-agent.js
Restart=always
RestartSec=1
StandardOutput=append:/tmp/vm-browser-agent.log
StandardError=append:/tmp/vm-browser-agent-error.log
Environment=HOME=/root
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
